# etc/profile.d/10umask.sh
# /etc/profile.d/**umask.sh
# set file mode creation mask, system wide. 
umask 022
