# etc/profile.d/90ashrc.sh
# /etc/profile.d/**ashrc.sh

# find and load ~/.ashrc
[ -f '~/.ashrc' ] && . ~/.ashrc
