# etc/profile.d/README
This directory should contain shell scripts configuring system-wide
environment on users' shells.

Files with the .sh extension found in this directory are evaluated by
Bourne-compatible shells (like ash, bash or zsh) when started as a
login shell.

These files are used each times you fire up a shell, you must not 
overload them too much as they setup default environment. Instead
you should offload definition to:

~/bin:/usr/local/bin:/usr/local/sbin
